# launchd/com.phoebe.chat-g-canva-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.phoebe.chat-g-canva-worker</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/phoebe/.codex/worktrees/e89a/課程招生 - 系統/scripts/chat-g-local-runner.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/phoebe/.codex/worktrees/e89a/課程招生 - 系統</string>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/phoebe/.codex/automations/chat-g-canva-proposals-automation/local-launchd.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/phoebe/.codex/automations/chat-g-canva-proposals-automation/local-launchd.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Applications/Codex.app/Contents/Resources</string>
    <key>ADMISSION_API_KEY</key>
    <string>admission-api-20260528-chat-a-trigger</string>
  </dict>
</dict>
</plist>
